# compose/server/supervisor/celery_beat.ini
[program:celery_beat]
command=celery -A server.celery_app beat -l INFO
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/celery_beat.out.log
stderr_logfile=/var/log/celery_beat.err.log
