# compose/server/supervisor/celery_dev.ini
[program:celery]
command=celery -A server.celery_app worker -B -l DEBUG
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/celery.out.log
stderr_logfile=/var/log/celery.err.log
