# compose/server/supervisor/celery_worker.ini
[program:celery_worker]
command=celery -A server.celery_app worker -l INFO
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/celery_worker.out.log
stderr_logfile=/var/log/celery_worker.err.log
