# compose/server/supervisor/gunicorn.ini
[program:gunicorn]
command=gunicorn server.app --access-logfile - --error-logfile - %(ENV_GUNICORN_ARGS)s
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/gunicorn.out.log
stderr_logfile=/var/log/gunicorn.err.log
